# chrome/root/etc/services.d/siyuan/run
#!/usr/bin/with-contenv bash

cd /opt/siyuan || exit

exec \
s6-setuidgid abc /opt/siyuan/kernel  -workspace=$WORKSPACE --accessAuthCode=$ACCESSAUTHCODE